# etc/netmon.conf
description "Manifold Map"
author      "Jordan Augé <jordan.auge@cisco.com>"

start on started  nfd
stop  on stopping nfd

respawn
respawn limit unlimited

script
  exec /usr/bin/netmon
end script
